# scripts/steamdeck/systemd/r4-f3-bridge.service
[Unit]
Description=DEPRECATED Wayfinder Aura R4 -> F3 bridge (legacy xdotool workaround)
# Superseded by wayfinder-trigger.service. Do not enable both.
Conflicts=wayfinder-trigger.service
After=graphical-session.target
PartOf=graphical-session.target
StartLimitIntervalSec=60
StartLimitBurst=5

[Service]
Type=simple
ExecStart=/home/deck/.local/bin/r4-f3-bridge.py
Restart=on-failure
RestartSec=3
StandardOutput=append:/tmp/r4-f3-bridge.log
StandardError=append:/tmp/r4-f3-bridge.log

[Install]
WantedBy=graphical-session.target
